# packages/pciauto.service
[Unit]
Description=pci

[Service]
ExecStart=/usr/bin/pcioff.sh

[Install]
WantedBy=multi-user.target
